# data/xdg-desktop-portal-gtk4.service.in
[Unit]
Description=Portal service (GTK4 implementation)
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.gtk4
ExecStart=@libexecdir@/xdg-desktop-portal-gtk4
